# roles/conditional_reboot/templates/conditional-reboot-check.service.j2
[Unit]
Description=conditional-reboot check

[Service]
ExecStart=/usr/local/bin/conditional-reboot -dry-run
Type=idle

[Install]
WantedBy=multi-user.target
