# playbooks/router/templates/services/metrics-emitter.j2
#!/bin/ksh

daemon="/etc/scripts/metrics-emitter"

. /etc/rc.d/rc.subr

pexp="/bin/sh ${daemon}"
rc_reload=NO
rc_bg=YES
rc_cmd $1
