# roles/ssh/templates/fetch-ssh-user-ca.timer.j2
[Unit]
Description=Run fetch-ssh-user-ca.service once a day

[Timer]
OnCalendar=daily
Persistent=true
Unit=fetch-ssh-user-ca.service

[Install]
WantedBy=timers.target
