# roles/minio/templates/minio.service.j2
[Unit]
Description=MinIO
Documentation=https://docs.min.io
Wants=network-online.target
After=network-online.target
AssertFileIsExecutable=/usr/local/bin/minio

# restart forever
StartLimitIntervalSec=0

[Service]
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectSystem=full
RestrictSUIDSGID=true
ReadWritePaths={{ minio_volumes }}
ReadOnlyPaths=/etc/minio.conf /home/{{ minio_user }}/.minio/certs/public.crt /home/{{ minio_user }}/.minio/certs/private.key

# --- Privilege restrictions ---
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
SecureBits=keep-caps-locked noroot noroot-locked

# --- Kernel / syscall hardening ---
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectHostname=true
ProtectClock=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=false
RemoveIPC=true

# Go binaries need a broader syscall surface than a typical C service
SystemCallArchitectures=native
SystemCallFilter=@system-service @resources
SystemCallFilter=~@debug @mount @cpu-emulation @obsolete @privileged @reboot @swap @raw-io @clock

# --- Network restrictions ---
IPAddressDeny=any
IPAddressAllow=127.0.0.0/8
IPAddressAllow=::1/128
IPAddressAllow=192.168.0.0/16

User={{ minio_user }}
Group={{ minio_group }}
ProtectProc=invisible

EnvironmentFile=/etc/minio.conf
ExecStartPre=/bin/bash -c "if [ -z \"${MINIO_VOLUMES}\" ]; then echo \"Variable MINIO_VOLUMES not set in /etc/minio.conf\"; exit 1; fi"
ExecStart=/usr/local/bin/minio server $MINIO_OPTS $MINIO_VOLUMES

# Let systemd restart this service always
Restart=always

# Specifies the maximum file descriptor number that can be opened by this process
LimitNOFILE=1048576

# Specifies the maximum number of threads this process can create
TasksMax=infinity

# Disable timeout logic and wait until process is stopped
TimeoutStopSec=infinity
SendSIGKILL=no

# restart forever
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target

# Built for ${project.name}-${project.version} (${project.name})
