# playbooks/vault-standby/templates/vault-cert-wrapper.service.j2
[Unit]
Description=Copy vault certificates

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vault-cert-wrapper.sh
