# roles/restic/templates/restic-init@.service.j2
[Unit]
Description=Init Restic repository

[Service]
Type=oneshot
ExecStart=/usr/local/bin/restic-upsert
EnvironmentFile=/etc/restic-%I.conf
