# playbooks/vault/templates/vault-cert-wrapper.timer.j2
[Unit]
Description=Run vault-cert-wrapper daily

[Timer]
OnBootSec=5min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
