# playbooks/router/templates/dhcpd-restart.service.j2
[Unit]
Description=Restart DHCPD service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart dhcpd.service
