# roles/unbound/templates/unbound_exporter.service.j2
[Unit]
Description=Unbound Prometheus Exporter
Documentation=https://github.com/letsencrypt/unbound_exporter
After=network.target
Wants=network.target

[Service]
Type=simple
User={{ unbound_exporter_user }}
Group={{ unbound_exporter_group }}
ExecStart=/usr/local/bin/unbound_exporter \
    -unbound.ca "" \
    -unbound.cert "" \
    -unbound.host "unix://{{ unbound_remote_control_socket_path }}" \
    -web.listen-address={{ unbound_exporter_listen_address }}

Restart=on-failure
RestartSec=60s
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/run

[Install]
WantedBy=multi-user.target
