# roles/smartctl_exporter/templates/smartctl-exporter.timer.j2
[Unit]
Description=Run smartmon script regularly

[Timer]
OnCalendar={{ smartctl_exporter_interval }}
Persistent=true

[Install]
WantedBy=timers.target
